# etc/supervisor/conf.d/condor.conf
[program:condor]
command=/bin/bash -c "/etc/init.d/condor start"
autostart=true
stdout_logfile= /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
